# meshcentral/rootfs/etc/s6-overlay/s6-rc.d/meshcentral/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: MeshCentral
# Starts MeshCentral
# ==============================================================================
bashio::log.info "Starting MeshCentral..."

export NODE_ENV=production

# Set log level if configured
declare log_level
log_level=$(bashio::config "log_level")
if bashio::var.has_value "${log_level}"; then
  export DEBUG="${log_level}"
fi

# Ensure data directories exist (defensive in case init step failed)
if ! bashio::fs.directory_exists "/share/meshcentral/meshcentral-data"; then
  mkdir -p /share/meshcentral/meshcentral-data
fi

if ! bashio::fs.directory_exists "/share/meshcentral/meshcentral-backups"; then
  mkdir -p /share/meshcentral/meshcentral-backups
fi

# Default domains configuration (always needed)
DEFAULT_DOMAINS='{"": {"title": "MeshCentral", "minify": true, "newAccounts": true, "myServer": {"Backup": true, "Restore": true}}}'

# Check if custom config is provided via addon option
declare meshcentral_config
meshcentral_config=$(bashio::config "meshcentral_config")

# Initialize or update config.json
if [ ! -s "/share/meshcentral/meshcentral-data/config.json" ] || bashio::var.has_value "${meshcentral_config}"; then
  if bashio::var.has_value "${meshcentral_config}"; then
    # Validate that it's not empty/whitespace/null
    if [[ "${meshcentral_config}" =~ ^[[:space:]]*$ ]] || [[ "${meshcentral_config}" == "null" ]]; then
      bashio::log.warning "meshcentral_config is empty or null, using default configuration"
      meshcentral_config=""
    else
      bashio::log.info "Using custom MeshCentral configuration from addon option..."
      # Validate JSON and ensure domains section exists using Node.js
      VALIDATED_CONFIG=$(node -e "
        const defaultDomains = JSON.parse(process.argv[2]);
        try {
          const config = JSON.parse(process.argv[1]);
          if (!config.domains || typeof config.domains !== 'object') {
            config.domains = defaultDomains;
          }
          console.log(JSON.stringify(config, null, 2));
        } catch (e) {
          console.error('Invalid JSON:', e.message);
          process.exit(1);
        }
      " "${meshcentral_config}" "${DEFAULT_DOMAINS}" 2>&1)
      VALIDATE_EXIT=$?
      
      if [ ${VALIDATE_EXIT} -eq 0 ] && [ -n "${VALIDATED_CONFIG}" ]; then
        echo "${VALIDATED_CONFIG}" > /share/meshcentral/meshcentral-data/config.json
        bashio::log.info "Custom configuration validated and written successfully"
      else
        bashio::log.error "Invalid JSON in meshcentral_config option"
        bashio::log.info "Falling back to default configuration..."
        meshcentral_config=""
      fi
    fi
  fi

  # Write default config if custom config wasn't used
  if ! bashio::var.has_value "${meshcentral_config}"; then
    bashio::log.info "Writing default MeshCentral configuration..."
    cat > /share/meshcentral/meshcentral-data/config.json <<EOF
{
  "settings": {
    "port": 44433,
    "redirPort": 44431,
    "exactPorts": true,
    "mpsPort": 44432,
    "relayPort": 453,
    "agentPort": 1234,
    "allowLoginToken": true,
    "allowFraming": true,
    "agentPing": 60,
    "agentPong": 60,
    "minify": true,
    "webRTC": true
  },
  "domains": ${DEFAULT_DOMAINS}
}
EOF
  fi
else
  bashio::log.info "Using existing MeshCentral configuration..."
fi

# Set data directory
export MESHCENTRAL_DATA="/share/meshcentral/meshcentral-data"
export MESHCENTRAL_FILES="/share/meshcentral/meshcentral-files"

cd /usr/local/lib/node_modules/meshcentral \
    || bashio::exit.nok "Could not change directory to MeshCentral"

exec node meshcentral.js --datapath /share/meshcentral/meshcentral-data
