# meshcentral/rootfs/etc/s6-overlay/s6-rc.d/meshcentral/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: MeshCentral
# Starts MeshCentral
# ==============================================================================
bashio::log.info "Starting MeshCentral..."

export NODE_ENV=production

# Set log level if configured
declare log_level
log_level=$(bashio::config "log_level")
if bashio::var.has_value "${log_level}"; then
  export DEBUG="${log_level}"
fi

# Ensure data directory exists (defensive in case init step failed)
if ! bashio::fs.directory_exists "/data/meshcentral/meshcentral-data"; then
  mkdir -p /data/meshcentral/meshcentral-data
fi

# Always write a known-good default config.json on startup
bashio::log.info "Writing default MeshCentral configuration..."
cat > /data/meshcentral/meshcentral-data/config.json <<EOF
{
  "settings": {
    "port": 80,
    "redirPort": 0,
    "exactPorts": true,
    "mpsPort": 44330,
    "relayPort": 453,
    "agentPort": 1234,
    "allowLoginToken": true,
    "allowFraming": true,
    "agentPing": 60,
    "agentPong": 60,
    "minify": true,
    "webRTC": true
  },
  "domains": {
    "": {
      "title": "MeshCentral",
      "minify": true,
      "newAccounts": true
    }
  }
}
EOF

# Set data directory
export MESHCENTRAL_DATA="/data/meshcentral/meshcentral-data"
export MESHCENTRAL_FILES="/data/meshcentral/meshcentral-files"

cd /usr/local/lib/node_modules/meshcentral \
    || bashio::exit.nok "Could not change directory to MeshCentral"

exec node meshcentral.js --datapath /data/meshcentral/meshcentral-data
