# meshcentral/rootfs/etc/s6-overlay/s6-rc.d/init-meshcentral/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: MeshCentral
# Creates data directories if they don't exist
# ==============================================================================

if ! bashio::fs.directory_exists "/share/meshcentral"; then
  mkdir -p /share/meshcentral
fi

if ! bashio::fs.directory_exists "/share/meshcentral/meshcentral-data"; then
  mkdir -p /share/meshcentral/meshcentral-data
fi

if ! bashio::fs.directory_exists "/share/meshcentral/meshcentral-files"; then
  mkdir -p /share/meshcentral/meshcentral-files
fi

if ! bashio::fs.directory_exists "/share/meshcentral/meshcentral-backups"; then
  mkdir -p /share/meshcentral/meshcentral-backups
fi

# Ensure proper permissions for MeshCentral data directories
chmod -R 755 /share/meshcentral 2>/dev/null || true
