# rootfs/etc/s6-overlay/s6-rc.d/init-meshcentral/up
#!/command/with-contenv bashio
# oneshot service for initializing MeshCentral data directories

